# mt5/ftmo/mt5-bridge-ftmo.service
[Unit]
Description=MT5 Bridge Proxy — FTMO (port 5056)
After=vncserver-ftmo.service

[Service]
Type=simple
WorkingDirectory=/home/tradebot/tradebots/mt5/common
Environment=WINEPREFIX=/home/tradebot/.wine64
Environment=WINEDEBUG=-all
Environment=DISPLAY=:1
Environment=MT5_BRIDGE_HOST=127.0.0.1
Environment=MT5_BRIDGE_PORT=5056
Environment=XDG_RUNTIME_DIR=/tmp/xdg-1000
Environment=TMPDIR=/tmp
ExecStartPre=/usr/bin/mkdir -p /tmp/xdg-1000
ExecStartPre=/usr/bin/chmod 700 /tmp/xdg-1000
ExecStart=/home/tradebot/tradebots/mt5/common/run_wine.sh mt5_bridge_proxy.py
StandardOutput=append:/home/tradebot/tradebots/mt5/ftmo/logs/bridge.log
StandardError=append:/home/tradebot/tradebots/mt5/ftmo/logs/bridge.log
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
